# cron/cookie_monster.cron
10 17 * * 0 echo "$(date) - Posting cookies from cron as user $(whoami)" >> /app/out/cron.log && PYTHONPATH=/app python3 /app/src/utils.py post_cookies
